# .config/systemd/user/nvidiasetfan.timer
[Unit]
Description=Run nvidiasetfan.service every 30s

[Timer]
OnCalendar=*:*:0/30

[Install]
WantedBy=timers.target