# .config/systemd/user/nvidiasetfan.service
[Unit]
Description=Run nvidiasetfan.sh

[Service]
Type=oneshot
ExecStart=%h/bin/nvidiasetfan.sh